# roles/nfs_automount/templates/nfs.automount
[Unit]
Description=Automount share

[Automount]
Where={{ nfs_automount_mount }}

[Install]
WantedBy=multi-user.target
